# qsnap@.timer
[Unit]
Description=qsnap BtrFS snapshot tool

[Timer]
OnCalendar=daily
Persistent=yes

[Install]
WantedBy=multi-user.target
