# qsnap@.service
[Unit]
Description=qsnap BtrFS snapshot tool

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mksnap.sh "%I"

[Install]
WantedBy=multi-user.target
